# projects/save_bathroom_info/bathroom.cron
*/5 12-22 * * * cd /home/yuanyuanzijin/code/dutsso-dev && python3 projects/save_bathroom_info/index.py > /dev/null 2>&1
